# deployment/branch-vm/systemd/sg-noc-nmap-discover.timer
[Unit]
Description=Run sg-noc-nmap-discover hourly

[Timer]
OnBootSec=10min
OnUnitActiveSec=1h
RandomizedDelaySec=5min
Persistent=true

[Install]
WantedBy=timers.target
